# .config/systemd/user/battery-notification.timer
[Unit]
Description=Run battery monitor notifications

[Timer]
OnBootSec=10
OnUnitActiveSec=10

[Install]
WantedBy=timers.target
