# .config/systemd/user/battery-notification.service
[Unit]
Description=Battery Notifications

[Service]
Type=simple
ExecStart=%h/.local/bin/wayde-battery-notification

[Install]
WantedBy=multi-user.target
