# .config/systemd/user/sway-battery-notify.service
[Unit]
Description=Battery Notifications

[Service]
Type=simple
ExecStart=%h/.local/bin/sway-battery-notify

[Install]
WantedBy=multi-user.target
