# .config/systemd/user/waybar-updates.service
[Unit]
Description=Check and cache software/firmware updates.

[Service]
Type=oneshot
ExecStart=%h/.local/bin/waybar-updates refresh all
